# deploy/roles/install_vmmaster/templates/upstart.conf
description 'Start vmmaster server'

start on startup
stop on shutdown

setuid vmmaster
chdir {{ app_path }}

respawn

script
    exec {{ app_path }}/.tox/bin/python manage.py runserver
end script