# deploy/roles/backend/templates/upstart_vmpool.conf
description 'Start vmpool server'

start on startup
stop on shutdown

setuid vmmaster
chdir {{ vmpool_path }}

respawn

script
    exec {{ vmpool_path }}/.tox/bin/python manage.py vmpool
end script